# model-api.service
[Unit]
Description=Model Selector API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/vibez/lab
ExecStart=/home/vibez/lab/rag-env/bin/python3 /home/vibez/lab/model-api.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
